# services/lispirc
#!/sbin/openrc-run

name="ii-lchat-worker-1"
description="IRC client worker #1"

command="/home/jose/src/sh/irc-connect"
command_args="(irc.libera.chat) <YOUR NICKSERV IDENT> <YOUR NICKSERV PASSWORD> <IRC FILES PATH> CHANNELS CHANNELS CHANNELS (eg. \#lisp)"
# for example:
#command_args="irc.libera.chat equwal 'PASSWORD WITH ESCAPES' /home/jose/scr/irc \#lisp \#commonlisp \#libreboot"
command_user="<YOUR USER>"
pidfile="/run/${RC_SVCNAME}.pid"
command_background="yes"

depend() {
    need net
    after logger
}

start_pre() {
    checkpath --directory --owner "$command_user" --mode 0755 /home/jose/scr/irc
}

stop() {
    if [ -f "$pidfile" ]; then
        local pid=$(cat "$pidfile")
        ebegin "Stopping $name"
        # Kill the main process and any child processes
        pkill -P "$pid" 2>/dev/null
        kill "$pid" 2>/dev/null
        # Wait a bit for graceful shutdown
        sleep 2
        # Force kill if still running
        if kill -0 "$pid" 2>/dev/null; then
            kill -KILL "$pid" 2>/dev/null
        fi
        # Clean up the pidfile
        rm -f "$pidfile"
        eend $?
    fi
}
